# .dotfiles/etc/rc.local
#!/usr/bin/env sh

chvt 2
plymouth --quit
